# scripts/octoprint.conf
[program:octoprint]
directory=/octoprint/
command=/octoprint/run --iknowwhatimdoing
numprocs=1
autostart=true
autorestart=true
user=root
